# ansible/files/clock-speed.timer
[Unit]
Description=Reset clock speed every 5 minutes

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
